# debian/usbbootgui.service
[Service]
Type=oneshot
ExecStart=/usr/bin/usbbootgui
